# Egor_Moysevich/0/files/crontab/list
* * * * * /home/vagrant/exercise-memcached.sh > /var/log/cron.log 2>&1
