# Egor_Moysevich/0/files/crontab/provision.sh
#!/bin/bash

crontab /tmp/crontab_list > /dev/null
echo 'crontab is ok'
